# src/nginx/tools/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log

[program:python_consumer]
command=python3 /rabbitmq_consumer.py
autostart=true
autorestart=true
stderr_logfile=/var/log/python_consumer.err.log
stdout_logfile=/var/log/python_consumer.out.log
